# common/supervisor/etc/php-fpm.ini
[program:php-fpm]
directory=/apps/srv/php56/sbin/
command=/apps/srv/php56/sbin/php-fpm
user=root
priority=999
autostart=true
autorestart=true
startsecs=3
startretries=2
stopsignal=KILL
stopwaitsecs=10
stdout_logfile=/apps/logs/supervisor/php-fpm.log
stderr_logfile=/apps/logs/supervisor/php-fpm.err

